# build/trustless.service
[Unit]
Description=Trustless daemon

[Service]
Type=simple
ExecStart=/usr/bin/trustless daemon

[Install]
WantedBy=default.target
