# scripts/trustless.service
[Unit]
Description=Trustless daemon

[Service]
Type=simple
ExecStart=@@@HOME@@@/bin/trustless daemon

[Install]
WantedBy=default.target
